# auditd_rules/tamandua.rules
# Tamandua EDR Audit Rules
# Generated for cross-platform ETW/auditd parity
# Performance mode: balanced
#
# This file provides comprehensive Linux audit rules that match the visibility
# of Windows ETW providers, enabling cross-platform threat detection.
#
# Installation:
#   1. Copy to /etc/audit/rules.d/tamandua.rules
#   2. Run: augenrules --load
#   3. Verify: auditctl -l | grep tamandua
#
# ETW Provider Mappings:
#   - Microsoft-Windows-Security-Auditing → Authentication & Authorization rules
#   - Microsoft-Windows-Sysmon → Process, File, Network rules
#   - Microsoft-Windows-PowerShell → Script execution monitoring
#   - Microsoft-Windows-WMI-Activity → Process monitoring
#   - Microsoft-Windows-Kernel-File → File operation rules
#   - Microsoft-Windows-Kernel-Network → Network operation rules
#   - Microsoft-Windows-Kernel-Process → Process creation rules

# Delete all existing rules
-D

# Set audit buffer size (balanced: 4096)
-b 4096

# Set failure mode (1=printk errors to dmesg)
-f 1

# Rate limit to prevent DoS (500 messages/sec)
-r 500

# ============================================================================
# Process Monitoring (ETW: Microsoft-Windows-Kernel-Process)
# ============================================================================

# Process creation (execve/execveat)
-a always,exit -F arch=b64 -S execve -k tamandua_process_create
-a always,exit -F arch=b32 -S execve -k tamandua_process_create
-a always,exit -F arch=b64 -S execveat -k tamandua_process_create
-a always,exit -F arch=b32 -S execveat -k tamandua_process_create

# Process injection (ptrace)
-a always,exit -F arch=b64 -S ptrace -k tamandua_process_inject
-a always,exit -F arch=b32 -S ptrace -k tamandua_process_inject

# Memory protection changes (mprotect - RWX transitions)
-a always,exit -F arch=b64 -S mprotect -k tamandua_mprotect
-a always,exit -F arch=b32 -S mprotect -k tamandua_mprotect

# ============================================================================
# File Operations (ETW: Microsoft-Windows-Kernel-File)
# ============================================================================

# File creation (open with O_CREAT flag)
-a always,exit -F arch=b64 -S open -S openat -S creat -F a1&0100 -k tamandua_file_create
-a always,exit -F arch=b32 -S open -S openat -S creat -F a1&0100 -k tamandua_file_create

# File deletion
-a always,exit -F arch=b64 -S unlink -S unlinkat -S rmdir -k tamandua_file_delete
-a always,exit -F arch=b32 -S unlink -S unlinkat -S rmdir -k tamandua_file_delete

# File rename
-a always,exit -F arch=b64 -S rename -S renameat -S renameat2 -k tamandua_file_rename
-a always,exit -F arch=b32 -S rename -S renameat -S renameat2 -k tamandua_file_rename

# File attribute changes
-a always,exit -F arch=b64 -S chmod -S fchmod -S fchmodat -k tamandua_file_attr
-a always,exit -F arch=b32 -S chmod -S fchmod -S fchmodat -k tamandua_file_attr
-a always,exit -F arch=b64 -S chown -S fchown -S lchown -S fchownat -k tamandua_file_attr
-a always,exit -F arch=b32 -S chown -S fchown -S lchown -S fchownat -k tamandua_file_attr

# Setuid/setgid file execution (privilege escalation)
-a always,exit -F arch=b64 -S execve -F perm=x -F auid>=1000 -F auid!=4294967295 -k tamandua_suid_exec
-a always,exit -F arch=b32 -S execve -F perm=x -F auid>=1000 -F auid!=4294967295 -k tamandua_suid_exec

# Sensitive directory monitoring (system binaries)
-w /bin -p wa -k tamandua_bin_changes
-w /sbin -p wa -k tamandua_bin_changes
-w /usr/bin -p wa -k tamandua_bin_changes
-w /usr/sbin -p wa -k tamandua_bin_changes
-w /usr/local/bin -p wa -k tamandua_bin_changes
-w /usr/local/sbin -p wa -k tamandua_bin_changes

# ============================================================================
# Network Operations (ETW: Microsoft-Windows-Kernel-Network)
# ============================================================================

# Socket creation
-a always,exit -F arch=b64 -S socket -k tamandua_socket_create
-a always,exit -F arch=b32 -S socket -k tamandua_socket_create

# Network connections (outbound)
-a always,exit -F arch=b64 -S connect -k tamandua_network_connect
-a always,exit -F arch=b32 -S connect -k tamandua_network_connect

# Network listening (inbound)
-a always,exit -F arch=b64 -S bind -k tamandua_network_bind
-a always,exit -F arch=b32 -S bind -k tamandua_network_bind
-a always,exit -F arch=b64 -S listen -k tamandua_network_listen
-a always,exit -F arch=b32 -S listen -k tamandua_network_listen
-a always,exit -F arch=b64 -S accept -S accept4 -k tamandua_network_accept
-a always,exit -F arch=b32 -S accept -S accept4 -k tamandua_network_accept

# ============================================================================
# Authentication & Authorization (ETW: Microsoft-Windows-Security-Auditing)
# ============================================================================

# Login/logout event logs
-w /var/log/lastlog -p wa -k tamandua_login
-w /var/log/faillog -p wa -k tamandua_login
-w /var/log/tallylog -p wa -k tamandua_login

# User and group management
-w /etc/passwd -p wa -k tamandua_identity
-w /etc/group -p wa -k tamandua_identity
-w /etc/shadow -p wa -k tamandua_identity
-w /etc/gshadow -p wa -k tamandua_identity
-w /etc/security/opasswd -p wa -k tamandua_identity

# Privilege escalation (sudo, su)
-w /etc/sudoers -p wa -k tamandua_priv_esc
-w /etc/sudoers.d/ -p wa -k tamandua_priv_esc
-a always,exit -F arch=b64 -S setuid -S setgid -S setreuid -S setregid -k tamandua_priv_esc
-a always,exit -F arch=b32 -S setuid -S setgid -S setreuid -S setregid -k tamandua_priv_esc

# PAM configuration changes
-w /etc/pam.d/ -p wa -k tamandua_pam
-w /etc/security/ -p wa -k tamandua_pam

# ============================================================================
# Privileged Operations (ETW: Microsoft-Windows-Threat-Intelligence)
# ============================================================================

# Kernel module operations (rootkit detection)
-a always,exit -F arch=b64 -S init_module -S finit_module -k tamandua_module_load
-a always,exit -F arch=b32 -S init_module -S finit_module -k tamandua_module_load
-a always,exit -F arch=b64 -S delete_module -k tamandua_module_unload
-a always,exit -F arch=b32 -S delete_module -k tamandua_module_unload
-w /etc/modprobe.d/ -p wa -k tamandua_module_config

# Capability changes (CAP_SYS_ADMIN, CAP_NET_ADMIN, etc.)
-a always,exit -F arch=b64 -S capset -k tamandua_capability
-a always,exit -F arch=b32 -S capset -k tamandua_capability

# LD_PRELOAD monitoring (library injection)
-w /etc/ld.so.preload -p wa -k tamandua_ld_preload

# System time changes (anti-forensics)
-a always,exit -F arch=b64 -S adjtimex -S settimeofday -S clock_settime -k tamandua_time_change
-a always,exit -F arch=b32 -S adjtimex -S settimeofday -S clock_settime -k tamandua_time_change
-w /etc/localtime -p wa -k tamandua_time_change

# ============================================================================
# Persistence Mechanisms (MITRE ATT&CK T1547)
# ============================================================================

# Cron job modifications
-w /etc/cron.d/ -p wa -k tamandua_cron
-w /etc/cron.daily/ -p wa -k tamandua_cron
-w /etc/cron.hourly/ -p wa -k tamandua_cron
-w /etc/cron.monthly/ -p wa -k tamandua_cron
-w /etc/cron.weekly/ -p wa -k tamandua_cron
-w /etc/crontab -p wa -k tamandua_cron
-w /var/spool/cron/ -p wa -k tamandua_cron

# Systemd unit modifications
-w /etc/systemd/system/ -p wa -k tamandua_systemd
-w /usr/lib/systemd/system/ -p wa -k tamandua_systemd
-w /lib/systemd/system/ -p wa -k tamandua_systemd

# Init script modifications
-w /etc/init.d/ -p wa -k tamandua_init
-w /etc/rc.local -p wa -k tamandua_init

# Shell profile modifications (bashrc, profile)
-w /etc/profile -p wa -k tamandua_profile
-w /etc/profile.d/ -p wa -k tamandua_profile
-w /etc/bash.bashrc -p wa -k tamandua_profile
-w /etc/bashrc -p wa -k tamandua_profile

# SSH authorized_keys modifications (lateral movement)
-a always,exit -F arch=b64 -S open -S openat -F dir=/home -F a1&0100 -F path~/.ssh/authorized_keys -k tamandua_ssh_keys
-a always,exit -F arch=b32 -S open -S openat -F dir=/home -F a1&0100 -F path~/.ssh/authorized_keys -k tamandua_ssh_keys
-w /root/.ssh/authorized_keys -p wa -k tamandua_ssh_keys

# ============================================================================
# Credential Access (MITRE ATT&CK T1003, T1552)
# ============================================================================

# Password hash file access
-a always,exit -F arch=b64 -S open -S openat -F path=/etc/shadow -k tamandua_credential_access
-a always,exit -F arch=b32 -S open -S openat -F path=/etc/shadow -k tamandua_credential_access

# SSH private key access
-a always,exit -F arch=b64 -S open -S openat -F dir=/home -F path~/.ssh/id_rsa -k tamandua_ssh_key_access
-a always,exit -F arch=b32 -S open -S openat -F dir=/home -F path~/.ssh/id_rsa -k tamandua_ssh_key_access
-a always,exit -F arch=b64 -S open -S openat -F dir=/home -F path~/.ssh/id_ed25519 -k tamandua_ssh_key_access
-a always,exit -F arch=b32 -S open -S openat -F dir=/home -F path~/.ssh/id_ed25519 -k tamandua_ssh_key_access

# Browser credential store access (Firefox, Chrome)
-a always,exit -F arch=b64 -S open -S openat -F dir=/home -F path~/.mozilla/firefox -k tamandua_browser_creds
-a always,exit -F arch=b32 -S open -S openat -F dir=/home -F path~/.mozilla/firefox -k tamandua_browser_creds
-a always,exit -F arch=b64 -S open -S openat -F dir=/home -F path~/.config/google-chrome -k tamandua_browser_creds
-a always,exit -F arch=b32 -S open -S openat -F dir=/home -F path~/.config/google-chrome -k tamandua_browser_creds

# ============================================================================
# Additional Security Monitoring
# ============================================================================

# Audit configuration tampering
-w /etc/audit/ -p wa -k tamandua_audit_config
-w /usr/sbin/auditctl -p x -k tamandua_audit_tools
-w /usr/sbin/auditd -p x -k tamandua_audit_tools

# Hosts file changes (DNS spoofing)
-w /etc/hosts -p wa -k tamandua_hosts_file

# Firewall rule changes (iptables, nftables)
-w /etc/iptables/ -p wa -k tamandua_firewall
-w /etc/nftables/ -p wa -k tamandua_firewall

# SELinux/AppArmor policy changes
-w /etc/selinux/ -p wa -k tamandua_selinux
-w /etc/apparmor/ -p wa -k tamandua_apparmor
-w /etc/apparmor.d/ -p wa -k tamandua_apparmor

# Package manager activity (yum, apt, dnf)
-w /var/log/yum.log -p wa -k tamandua_software
-w /var/log/apt/ -p wa -k tamandua_software
-w /var/log/dpkg.log -p wa -k tamandua_software

# ============================================================================
# Performance Notes
# ============================================================================
#
# These rules are optimized for balanced performance (~5-10% CPU overhead):
# - Buffer size: 4096 events
# - Rate limit: 500 events/second
# - Excludes high-frequency innocuous operations
# - Uses arch-specific filters (b64/b32)
# - Leverages key-based event tagging for fast filtering
#
# For aggressive mode: increase buffer (-b 8192) and remove rate limit
# For lightweight mode: reduce buffer (-b 1024) and rate limit (-r 100)
#
# ============================================================================
# End of Tamandua EDR Audit Rules
# ============================================================================
